# combined/sv/nginx/finish
#!/bin/sh -e

if [ "${SNIKKET_TWEAK_WEB_PROXY:-1}" = "0" ]; then
	# tell s6 it's not supposed to be running
	exit 125
fi
