# combined/sv/nginx/run
#!/bin/sh -e

exec /usr/sbin/nginx -g "daemon off;"
