# ansible/files/refresh-certs.cron
#!/bin/sh

prosodyctl --root cert import /snikket/letsencrypt/live
