# ansible/files/supervisor-mailhog.conf
[program:mailhog]
command=/usr/local/bin/mailhog -auth-file=/etc/mailhog-auth
autorestart=true
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
umask=002
