# ansible/files/supervisor-coturn.conf
[program:coturn]
command=start-coturn.sh
startsecs=0
autorestart=unexpected
exitcodes=0
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
umask=002
