# ansible/files/supervisor-certmon.conf
[program:cert-monitor]
command=cert-monitor.sh
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
umask=002
