# combined/sv/anacron/run
#!/bin/sh -e

while /usr/sbin/anacron -d -n; do
	sleep 3600;
done
